# amis/base/root/etc/systemd/system/consul-registration.timer
[Unit]
Description="Run script after 5 minutes of boot"

[Timer]
OnBootSec=5min

[Install]
WantedBy=default.target